# infrastructure/server-hardening/systemd/empire-monitor.service
[Unit]
Description=Empire Server Health Monitor
After=network.target docker.service

[Service]
Type=oneshot
ExecStart=/opt/monitor/watchdog.sh
StandardOutput=journal
StandardError=journal
